# INT26-6/systemd/disk-monitor.service
[Unit]
Description=Disk Monitor Service

[Service]
Type=oneshot
User=monitoruser
Group=monitoruser
ExecStart=/usr/local/src/disk_monitor.sh
