# INT26-6/systemd/disk-monitor.timer
[Unit]
Description=Disk Monitor Timer

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
