# INT26-6/systemd/ram-monitor.timer
[Unit]
Description=RAM Monitor Timer

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
