# INT26-6/systemd/ram-monitor.service
[Unit]
Description=RAM Monitor Service

[Service]
Type=oneshot
User=monitoruser
Group=monitoruser
ExecStart=/usr/local/src/ram_monitor.sh
